# scripts/docker/supervisor-catmaid.conf
[program:uwsgi-catmaid]
command = /home/scripts/docker/start-catmaid.sh
stopsignal = INT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
