# scripts/docker/supervisor-catmaid.conf
[program:uwsgi-catmaid]
command = /home/scripts/docker/start-catmaid.sh
