# scripts/docker/supervisor-catmaid.conf
[program:postgresql]
command = service postgresql start
startsecs = 0

[program:uwsgi-catmaid]
command = /home/scripts/docker/start-catmaid.sh

[program:nginx]
command = /usr/sbin/nginx
