# scripts/docker/supervisor-catmaid.conf
[program:postgresql]
command = service postgresql start
startsecs = 0

[program:uwsgi-catmaid]
command = /usr/local/bin/uwsgi --ini /home/scripts/docker/uwsgi-catmaid.ini

[program:nginx]
command = /usr/sbin/nginx
